# docker/supervisor.d/app.conf
[program:app]
process_name=app
command=php -S 0.0.0.0:8000 -t ./public/
autostart=true
autorestart=true
redirect_stderr=true
