# docker/supervisor.d/worker.conf
[program:worker]
process_name=worker
command=php /var/www/html/console.php queue:work
autostart=true
autorestart=true
redirect_stderr=true
